# systemd/mailman.timer
[Unit]
Description=Mailman

[Timer]
OnUnitActiveSec=30s
OnBootSec=30s

[Install]
WantedBy=timers.target