# systemd/mailman.service
[Unit]
Description=Mailman Service

[Service]
Type=oneshot
ExecStart=/usr/bin/mailman --systemd 8.8.8.8